# s6/confd/run
#!/bin/bash

exec confd -watch --prefix=$PREFIX
